# systemd/user/emacs-daemon.service
[Unit]
Description=emacs daemon

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/emacs --daemon

[Install]
WantedBy=default.target
