# systemd/user/i3.target
[Unit]
Description=i3
After=xinitrc.target
Wants=xinitrc.target
Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
