# systemd/user/offlineimap.service
[Unit]
Description=Offlineimap

[Service]
ExecStart=/usr/bin/offlineimap
Restart=always

[Install]
WantedBy=default.target

