# systemd/user/test.service
[Unit]
Description=test
[Service]
ExecStart=/bin/sh -c 'dmesg | tac'
