# systemd/user/xinitrc.target
[Unit]
Description=xinitrc stuff
Requires=xorg.target
After=xorg.target

